# roles/admin/templates/opsgenie/opsgenie_heartbeat.timer.j2
[Unit]
Description=Timer for OpsGenie Heartbeat {{ item.name }}
Requires=opsgenie-{{ item.name }}-heartbeat.service

[Timer]
OnCalendar=*:0/{{ item.interval }}

[Install]
WantedBy=timers.target
